# debian/museobutton.postinst
#!/bin/bash

if [ -x "/etc/init.d/museobutton" ]; then
	update-rc.d museobutton defaults >/dev/null
	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
		invoke-rc.d museobutton start || exit 0
	else
		/etc/init.d/museobutton start || exit 0
	fi
fi

