# debian/museobutton.cron.d

# start the daemon if not already started
* * * * *	root	/usr/sbin/invoke-rc.d museobutton start
