# debian/museobutton.cron.d

# start the daemon if not already started
* * * * *	root	invoke-rc.d museobutton start
